# deploy/systemd/at-yaris-race-freeze.timer
[Unit]
Description=Check Race-Level Final Prediction Windows Every Two Minutes

[Timer]
OnCalendar=*-*-* 09..23:00/2:00 Europe/Istanbul
AccuracySec=1s
Persistent=true
Unit=at-yaris-race-freeze.service

[Install]
WantedBy=timers.target
